# common/recipes-core/udev/udev-extraconf/60-drm-hotplug-mode.rules
ACTION!="remove", SUBSYSTEM=="drm", SUBSYSTEMS=="pci|usb|platform", RUN+="/etc/udev/scripts/hotplug-display-@MULTIDISPLAYMODE.sh"
